# squashfs-root/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

#
rc_file="/data/etc/rc.local"
[ -f "$rc_file" ] || cp /etc/rc.template $rc_file

. "$rc_file"

services_file="/data/etc/services"
[ -f "$rc_file" ] || cp /etc/services.template $services_file

exit 0
